# foo.cron
10 14 * * 1 /home/pi/Downloads/Mikrotik-Blacklist/filesplitterv2.sh
